# wildlife-tracker.service
[Unit]
Description=Wildlife PTZ Camera Tracker backend
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=gojian
Group=gojian
WorkingDirectory=/home/gojian/Code/Wildlife_PTZ_Camera_Tracker
ExecStart=/home/gojian/Code/Wildlife_PTZ_Camera_Tracker/.venv/bin/uvicorn \
    backend.main:app \
    --host 0.0.0.0 \
    --port 9090
Restart=on-failure
RestartSec=5

Environment=PYTHONUNBUFFERED=1
# Jetson Orin only: torch 2.5.0a0 (NVIDIA JetPack build) links against libcusparseLt.so.0
# at runtime.  JetPack does not ship this library; it lives in the Python venv instead.
# On desktop CUDA systems the path below does not exist — the setting is harmless.
Environment=LD_LIBRARY_PATH=/home/gojian/Code/Wildlife_PTZ_Camera_Tracker/.venv/lib/python3.10/site-packages/nvidia/cusparselt/lib

# Give the GPU time to initialise on boot before the first request hits
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
